# systemd/larperdetector.service.example
[Unit]
Description=VATSIM Watchlist Bot (example)
After=network.target

[Service]
# Adjust User/Group and WorkingDirectory to your install path
User=bot
Group=bot
WorkingDirectory=/opt/vatsim-watchlist-bot

# If you're using a virtualenv, set ExecStart to the venv python path, e.g.
# ExecStart=/opt/vatsim-watchlist-bot/.venv/bin/python /opt/vatsim-watchlist-bot/bot.py
ExecStart=/usr/bin/python3 /opt/vatsim-watchlist-bot/bot.py

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
